# systemd/system/snap-google\x2dcloud\x2dsdk-148.mount
[Unit]
Description=Mount unit for google-cloud-sdk, revision 148
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/google-cloud-sdk_148.snap
Where=/snap/google-cloud-sdk/148
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
